# etc/init.d/tiny-cloud-early
#!/sbin/openrc-run
# vim:set ft=sh noet ts=4:

description="Tiny Cloud Bootstrap - early phase"

depend() {
	before mdev
}

start() {
	source /lib/tiny-cloud/init-early

	is_bootstrapped && return 0

	ebegin "Expanding Root Volume/Partition"
	expand_root
	eend $?

	# TODO: _setup mdev things, if applicable
}