# dist/openrc/tiny-cloud-main
#!/sbin/openrc-run
# vim:set ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - main phase"

depend() {
	need net
	before sshd
}

start() {
	ebegin "Tiny Cloud - main phase"
	tiny-cloud main
	eend $?
}
