# dist/openrc/tiny-cloud-early
#!/sbin/openrc-run
# vim:set filetype=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - early phase"

depend() {
	need net
	before tiny-cloud-main
}

start() {
	ebegin "Tiny Cloud - early phase"
	tiny-cloud early
	eend $?
}
