# dist/openrc/tiny-cloud
#!/sbin/openrc-run
# vim:set ts=8 noet ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - main phase"
extra_commands="complete incomplete"

depend() {
	need net
	before sshd
}

start() {
	ebegin "Tiny Cloud - main phase"
	tiny-cloud main
	eend $?
}

# allow setting / unsetting of bootstrapped state
complete() { tiny-cloud --bootstrap complete ; }
incomplete() { tiny-cloud --bootstrap incomplete ; }
