# dist/openrc/tiny-cloud-final
#!/sbin/openrc-run
# vim:set ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - final phase"

depend() {
	after *
	provide cloud-final
}

start() {
	ebegin "Tiny Cloud - final phase"
	tiny-cloud final
	eend $?
}
