# dist/openrc/tiny-cloud-early
#!/sbin/openrc-run
# vim:set ts=8 noet ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - early phase"

depend() {
	after root logger
	before net
}

start() {
	ebegin "Tiny Cloud - early phase"
	tiny-cloud early
	eend $?
}
