# dist/openrc/tiny-cloud-boot
#!/sbin/openrc-run
# vim:set ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - boot phase"

depend() {
	after root logger
	before net
}

start() {
	ebegin "Tiny Cloud - boot phase"
	tiny-cloud boot
	eend $?
}
