# dist/openrc/tiny-cloud-net
#!/sbin/openrc-run
# vim:set ts=8 noet ft=sh:
# shellcheck shell=sh

description="Tiny Cloud Bootstrap - net phase"
depend() {
	need net
	before tiny-cloud
}

start() {
	ebegin "Tiny Cloud - net"
	tiny-cloud net
	eend $?
}
